# home/config/systemd/user/mc-backup.service
[Unit]
Description=Snapshot GDLauncher Minecraft Instances

[Service]
Type=oneshot
ExecStart=%h/.config/dotfiles-scripts/minecraft_backup.sh
